# supervisord-conf/genesisapi.ini
[program:genesisapi]

command = sh /home/genesis/src/genesis-gunicorn.sh
directory = /home/genesis/src/
user = genesis
autosart = true 
autorestart = true 
stdout_logfile=/home/genesis/log/genesisapi.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=5
stderr_logfile=/home/genesis/log/genesisapi.log
stderr_logfile_maxbytes=5MB
stderr_logfile_backups=5